# Server/ets/conf/Ubuntu upstart/ets-server.conf

description "ETS server"
author      "EveryAware"

start on filesystem and started networking and redis-server-started
stop on shutdown

respawn             # restart when job dies
respawn limit 5 60  # give up restart after 5 respawns in 60 seconds

env NODE_ENV=production
env NODE_BIN=/usr/local/bin/node

env ETS_LOG=/var/log/ets.log
env ETS_ROOT=/var/www/yourPath/ets/node
env ETS_PID=/var/run/ets.pid
env ETS_USER=www-data

script
  exec start-stop-daemon --start --chdir ${ETS_ROOT} --make-pidfile --pidfile ${ETS_PID} -c ${ETS_USER} --exec ${NODE_BIN} server.js >> ${ETS_LOG} 2>&1
end script

pre-start script
  # Date format same as (new Date()).toISOString() for consistency
  echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> ${ETS_LOG}
end script

pre-stop script
  rm -f ${ETS_PID}
  echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> ${ETS_LOG}
end script

